# src/supervisord.conf
[supervisord]
nodaemon=true

[program:incron]
command=/usr/sbin/incrond

[program:rsync]
command=/usr/bin/rsync --no-detach --daemon --config /etc/rsyncd.conf
